# crontab
# Run metering processor every 5 minutes (with lock to prevent overlapping)
*/5 * * * * cd /app && flock -n /tmp/metering_processor.lock python3 src/metering_processor.py >> /var/log/cron.log 2>&1

# Empty line required at end of crontab
