# crontab
# Run metering processor every 5 minutes
*/5 * * * * cd /app && python3 src/metering_processor.py >> /var/log/cron.log 2>&1

# Empty line required at end of crontab
